# templates/chronos.service.erb
[Unit]
Description=Chronos
After=network.target
Wants=network.target

[Service]
EnvironmentFile=-/etc/sysconfig/chronos
ExecStart=/usr/bin/sh -c "exec /usr/bin/chronos ${args}"
Restart=always
RestartSec=20

[Install]
WantedBy=multi-user.target
